# docker/supervisord.conf
; Supervisor: mantem Apache + worker de fila + scheduler vivos no mesmo container.
; Usado como entrypoint pelo scripts/start-container.sh (producao / EasyPanel).

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; Servidor web (processo principal do container).
[program:apache]
command=apache2-foreground
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Worker da fila: processa SyncAgendaEventToCalendarsJob, ProcessCalendarWebhookJob, etc.
; --max-time recicla o processo de hora em hora para liberar memoria (o Supervisor reinicia).
[program:queue-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
stopwaitsecs=3600
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Scheduler: dispara as tarefas de routes/console.php (renovacao do watch a cada 6h,
; lembretes a cada 5 min, DataJud diario). schedule:work roda em primeiro plano e
; aciona o que estiver vencido a cada minuto.
[program:scheduler]
command=php /var/www/html/artisan schedule:work
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
